# scripts/podman/openclaw.container.in
# FoxClaw gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into foxclaw's ~/.config/containers/systemd/
# {{FOXCLAW_HOME}} is replaced at install time.

[Unit]
Description=FoxClaw gateway (rootless Podman)

[Container]
Image=foxclaw:local
ContainerName=foxclaw
UserNS=keep-id
# Keep container UID/GID aligned with the invoking user so mounted config is readable.
User=%U:%G
Volume={{FOXCLAW_HOME}}/.foxclaw:/home/node/.foxclaw:Z
EnvironmentFile={{FOXCLAW_HOME}}/.foxclaw/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
